# home-assistant-addon/tests/supervisord-test.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
loglevel=info

[program:nginx]
command=nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=0

[program:mastra-server]
command=node --import=./mcp/.mastra/output/instrumentation.mjs ./mcp/.mastra/output/index.mjs
directory=/workspace
environment=PORT="%(ENV_MASTRA_SERVER_PORT)s",HOST="0.0.0.0"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=1

[program:mcp-server]
command=/workspace/mcp/mcp-server
directory=/workspace
environment=PORT="%(ENV_MCP_SERVER_PORT)s",HOST="0.0.0.0"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=2
