# home-assistant-addon/tests/supervisord-test.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
loglevel=info

[program:nginx]
command=nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=0

[program:mastra-dev]
command=bunx mastra dev --dir mcp/mastra --root .
directory=/workspace
environment=PORT="%(ENV_MASTRA_UI_PORT)s",HOST="0.0.0.0"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=1

[program:mcp-server]
command=bun run mcp/mastra/mcp-server.ts
directory=/workspace
environment=PORT="%(ENV_MCP_SERVER_PORT)s",HOST="0.0.0.0"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=2
